# proxy/rootfs/etc/s6-overlay/s6-rc.d/lunar-engine/run
#!/command/with-contenv sh

exec lunar_engine
