# proxy/rootfs/etc/s6-overlay/s6-rc.d/lunar-engine/up
#!/bin/sh

s6-envdir -Lfn -- /run/s6/container_environment

exec lunar_engine
