# proxy/rootfs/etc/s6-overlay/s6-rc.d/lunar-cron/up
foreground { /etc/s6-overlay/s6-rc.d/lunar-cron/setup-cron.sh }

sudo /usr/sbin/service cron start