# proxy/rootfs/etc/s6-overlay/s6-rc.d/fluent-bit/run
#!/command/execlineb -P

s6-envdir -Lfn -- /run/s6/container_environment
s6-envdir /var/run/lunar_env

trap
{
  cont
  {
    exit 0
  }
  term
  {
    /usr/bin/stop_services fluent-bit
  }
}

exec fluent-bit -c /etc/fluent-bit/fluent-bit.conf -e /etc/fluent-bit/plugin/output_aggregation.so --enable-hot-reload
