# proxy/rootfs/etc/s6-overlay/s6-rc.d/async-service/run
#!/command/execlineb -P

s6-envdir -Lfn -- /run/s6/container_environment


trap
{
  term
  {
    /usr/bin/stop_services async-service
  }
}

exec async-service
