# proxy/rootfs/etc/s6-overlay/s6-rc.d/logrotate/run
#!/command/execlineb -P

exec /usr/bin/logrotate
