# proxy/rootfs/etc/s6-overlay/s6-rc.d/haproxy/run
#!/command/execlineb -P

s6-envdir -Lfn -- /run/s6/container_environment

trap
{
  term
  {
    /usr/bin/stop_services haproxy
  }
}

foreground
{
  exec /usr/bin/check_for_tls
}

exec haproxy -W -db -f /etc/haproxy/haproxy.cfg
