# proxy/rootfs/etc/s6-overlay/s6-rc.d/haproxy/run
#!/command/execlineb -P

s6-envdir -Lfn -- /run/s6/container_environment

trap
{
# TERM broadcast is override in order to shut down gracefully

  term
  {
    /usr/bin/stop_services haproxy
  }
}

haproxy -W -db -f /etc/haproxy/haproxy.cfg
